# images/app/etc/service/apache2/run
#!/bin/sh
exec 2>&1
cd /opt/logjam/app
exec passenger start --port 80 --log-file=/dev/stdout
