# images/app/etc/service/logjam/run
#!/bin/sh
exec 2>&1
exec chpst -ulogjam runsvdir /opt/logjam/app/service .........................................
