# images/demo/service/mongodb/run
#!/bin/sh
mkdir -p /data/db
exec 2>&1
exec mongod
