# images/demo/etc/service/apache2/run
#!/bin/sh
exec 2>&1
exec /usr/sbin/apache2ctl -D FOREGROUND
